# .ci/phpunit-environment.conf
[program:php-webserver]
command=php ; php-webserver -S localhost:8080
directory=. ; php-webserver
autostart=false ; php-webserver
